# deploy/huduglue-monitor.service
[Unit]
Description=HuduGlue Website Monitor Check
After=network.target

[Service]
Type=oneshot
User=administrator
Group=administrator
WorkingDirectory=/home/administrator
Environment="PATH=/home/administrator/venv/bin"
EnvironmentFile=/home/administrator/.env
ExecStart=/home/administrator/venv/bin/python manage.py check_websites
StandardOutput=journal
StandardError=journal
